# syncthing/.config/systemd/user/syncthing.service
[Unit]
Description=Syncthing

[Service]
Type=simple
ExecStart=syncthing

[Install]
WantedBy=default.target
